# liveBusDelayHeatMap/crontab
0      * * * *  root    sh /app/getDelayData_full.sh >> /app/cron_full_latest.log 2>&1
10-50/10 * * * *  root    sh /app/getDelayData_incremental.sh >> /app/cron_inc_latest.log 2>&1
